# etc/crontab
10 5 * * * /root/rocketchat/cron_job.sh 
